# app/bokeh.conf
[program:myserver]
command=bokeh serve --host '*'

[supervisord]
logfile=/home/ubuntu/supervisord_bokeh.log
loglevel=debug
user=root
